# supervisord.conf
# supervisord.conf
# Supervisor is the process manager that starts and monitors all services
# inside the container. It replaces a traditional init system (like systemd)
# which doesn't run inside Docker.

[supervisord]
nodaemon=true                                      # Run in the foreground so Docker can track the process
user=root                                          # Run supervisor itself as root so it can manage system services
logfile=/var/log/supervisor/supervisord.log        # Main supervisor log
pidfile=/var/run/supervisord.pid                   # PID file location

# ── 1. D-Bus ─────────────────────────────────────────────────────────────────
# D-Bus provides inter-process communication used by the desktop and apps.
# Must start before the display server and desktop session.
[program:dbus]
command=/usr/bin/dbus-daemon --system --nofork --nopidfile
autorestart=true
priority=5                                         # Lowest number = starts first
stdout_logfile=/var/log/supervisor/dbus.log
stderr_logfile=/var/log/supervisor/dbus_err.log

# ── 2. Xvnc (VNC Display Server) ─────────────────────────────────────────────
# Xvnc is a virtual X display server with a built-in VNC server.
# It creates display :1 at 1920x1080, 24-bit colour, on port 5901.
#   -localhost no    → allow connections from outside the container
#   -dpi 96          → standard screen DPI
#   -ac              → disable access control (safe inside Docker)
#   -rfbauth         → path to the VNC password file
[program:xvnc]
command=/usr/bin/Xvnc :1 -geometry 1920x1080 -depth 24 -rfbport 5901 -localhost no -dpi 96 -ac -rfbauth /home/user/.vnc/passwd
autorestart=true
priority=10                                        # Starts after dbus
stdout_logfile=/var/log/supervisor/xvnc.log
stderr_logfile=/var/log/supervisor/xvnc_err.log

# ── 3. Cinnamon Desktop Session ───────────────────────────────────────────────
# Launches the Cinnamon desktop environment on the Xvnc display.
# Environment variables set here ensure the session uses software rendering
# (no real GPU) and targets display :1.
[program:cinnamon]
command=/bin/bash /usr/local/bin/start-cinnamon.sh
autorestart=true
priority=20                                        # Starts after Xvnc is ready
environment=HOME="/home/user",USER="user",DISPLAY=":1",LIBGL_ALWAYS_SOFTWARE="1",CLUTTER_BACKEND="x11"
stdout_logfile=/var/log/supervisor/cinnamon.log
stderr_logfile=/var/log/supervisor/cinnamon_err.log

# ── 4. noVNC (Browser-based VNC client) ──────────────────────────────────────
# websockify proxies WebSocket connections (port 6080) to the VNC server (5901),
# allowing access via a web browser without a dedicated VNC client.
# Note: port 6080 is not exposed in docker-compose.yml by default.
[program:novnc]
command=/usr/bin/websockify --web /usr/share/novnc/ 6080 localhost:5901
autorestart=true
priority=30                                        # Starts last
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc_err.log
